# Docker/supervisor.conf
[supervisord]
nodaemon=true

[program:mongo]
command=mongod

[program:scheduler]
command=java -jar Scheduler.jar